# etc/rc.local
#!/bin/ash -e

/lib/rtl/init.sh &
#/lib/mmc/init.sh &
#/lib/mmc/p28.sh &

exit 0
